# ops/systemd/tuwayki-surface@.service
[Unit]
Description=TUWAYKI Reflex Surface (%i)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/tuwayki
EnvironmentFile=/etc/tuwayki/%i.env
ExecStart=/bin/bash -lc '/opt/tuwayki/.venv/bin/reflex run --env prod --single-port --frontend-port ${REFLEX_PORT}'
Restart=always
RestartSec=3
TimeoutStartSec=120
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
